# package/home-webserver-init-script
#!/sbin/openrc-run

name="home-webserver"
description="Home Web Server"
command="/usr/local/bin/home-webserver"
command_args="--config /etc/home-webserver/config.toml"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath -f -m 0644 -o root:root "$output_log"
    checkpath -d -m 0755 -o root:root /etc/home-webserver
    
    if [ ! -f /etc/home-webserver/users.txt ]; then
        eerror "Users file not found at /etc/home-webserver/users.txt"
        return 1
    fi
    
    if [ ! -f /etc/home-webserver/config.toml ]; then
        eerror "Config file not found at /etc/home-webserver/config.toml"
        return 1
    fi
}