# package/openrc-init
#!/sbin/openrc-run

name="dashboard-web"
description="Web Dashboard"
command="/usr/local/bin/dashboard-web"
command_args="--config /etc/dashboard-web/config.toml"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
    need net
    after firewall
}

start_pre() {
    checkpath -d -m 0755 -o root:root /etc/dashboard-web
    
    if [ ! -f /etc/dashboard-web/config.toml ]; then
        eerror "Config file not found at /etc/dashboard-web/config.toml"
        return 1
    fi
}

reload() {
    ebegin "Reloading ${RC_SVCNAME} configuration"
    if [ -f "${pidfile}" ]; then
        local pid=$(cat "${pidfile}")
        if [ -n "${pid}" ] && kill -0 ${pid} 2>/dev/null; then
            kill -HUP ${pid}
            eend $? "Failed to reload ${RC_SVCNAME}"
        else
            eend 1 "Process ${pid} not running"
        fi
    else
        eend 1 "Pidfile ${pidfile} not found"
    fi
}